# docker/allinone/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql-err.log

[program:redis]
command=/usr/bin/redis-server --maxmemory 256mb --maxmemory-policy allkeys-lru --appendonly yes --dir /data/redis
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-err.log

[program:backend-web]
command=python /app/backend-web/main.py
directory=/app/backend-web
autostart=true
autorestart=true
priority=20
startsecs=10
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/backend-web.log
stderr_logfile=/var/log/supervisor/backend-web-err.log

[program:websocket]
command=python /app/websocket/main.py
directory=/app/websocket
autostart=true
autorestart=true
priority=30
startsecs=10
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/websocket.log
stderr_logfile=/var/log/supervisor/websocket-err.log

[program:scheduler]
command=python /app/scheduler/main.py
directory=/app/scheduler
autostart=true
autorestart=true
priority=30
startsecs=10
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/scheduler.log
stderr_logfile=/var/log/supervisor/scheduler-err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-err.log
